# docker/cron/crontab
*/5 * * * * root /usr/local/bin/script.sh >> /var/log/cron.log 2>&1
